# supervisor-app.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=uwsgi --ini /app/uwsgi.ini

[program:nginx]
command=service nginx start
